# .config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon
After=dbus.socket

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target
